# keyringer/assets/keyringer.timer
[Unit]
Description=Run KeyRinger Service Daily

[Timer]
OnCalendar=daily
Unit=keyringer.service
Persistent=true

[Install]
WantedBy=timers.target